# ROOT/etc/crontab
30 3 * * * su server sh /home/server/app/check.sh
